# docker/etc/service/registration_server/run
#!/bin/sh

export RUST_LOG=${RUST_LOG:-info}
exec 2>&1
exec /home/user/registration_server/target/release/main --config-file=/home/user/config/config.toml
