# docker/etc/service/pagekite/run
#!/bin/sh

exec 2>&1
exec /usr/local/bin/pagekite.py --optfile=/home/user/config/pagekite.conf
