# docker/etc/cron.weekly/geoipupdate
#!/bin/sh

if grep -q ^AccountID /etc/GeoIP.conf; then
    /usr/bin/geoipupdate
fi
